# data/qnguyenle/bge-attackqa-retriever/model.yaml
_id: 69ef150c5624f54038cda158
id: qnguyenle/bge-attackqa-retriever
private: false
pipeline_tag: sentence-similarity
library_name: sentence-transformers
tags:
- sentence-transformers
- safetensors
- bert
- sentence-similarity
- feature-extraction
- generated_from_trainer
- dataset_size:22801
- loss:MultipleNegativesRankingLoss
- arxiv:1908.10084
- arxiv:1807.03748
- base_model:BAAI/bge-base-en-v1.5
- base_model:finetune:BAAI/bge-base-en-v1.5
- text-embeddings-inference
- endpoints_compatible
- region:us
downloads: 3
likes: 0
modelId: qnguyenle/bge-attackqa-retriever
author: qnguyenle
sha: cd9cbbcaa6243c1c47c2ec58dbb6f679b5e95978
lastModified: '2026-04-27T07:50:09.000Z'
gated: false
disabled: false
mask_token: '[MASK]'
widgetData:
- source_sentence: 'Represent this sentence for searching relevant passages: How does attack software ''S0030: Carbanak''
    use attack technique ''T1027: Obfuscated Files or Information''?'
  sentences:
  - 'How attack software ''S0363: Empire'' uses attack technique ''T1558.002: Silver Ticket'':

    Empire can leverage its implementation of Mimikatz to obtain and use silver tickets.'
  - 'How attack software ''S0664: Pandora'' uses attack technique ''T1027: Obfuscated Files or Information'':

    Pandora has the ability to compress stings with QuickLZ.'
  - 'How attack software ''S0030: Carbanak'' uses attack technique ''T1027: Obfuscated Files or Information'':

    Carbanak encrypts strings to make analysis more difficult.'
- source_sentence: 'Represent this sentence for searching relevant passages: What type of organizations has MegaCortex mainly
    targeted?'
  sentences:
  - 'How attack software ''S0534: Bazar'' uses attack technique ''T1070.009: Clear Persistence'':

    Bazar''s loader can delete scheduled tasks created by a previous instance of the malware.'
  - 'Description of attack software ''S0576: MegaCortex'':

    MegaCortex is ransomware that first appeared in May 2019.  MegaCortex has mainly targeted industrial organizations.'
  - 'How attack software ''S0093: Backdoor.Oldrea'' uses attack technique ''T1057: Process Discovery'':

    Backdoor.Oldrea collects information about running processes.'
- source_sentence: 'Represent this sentence for searching relevant passages: Describe attack group ''G0024: Putter Panda'''
  sentences:
  - 'How data component ''Script Execution'' can be used to detect attack technique ''T1005: Data from Local System'':

    Monitor for any attempts to enable scripts running on a system would be considered suspicious. If scripts are not commonly
    used on a system, but enabled, scripts running out of cycle from patching or other administrator functions are suspicious.
    Scripts should be captured from the file system when possible to determine their actions and intent.  Data may also be
    acquired through Windows system management tools such as Windows Management Instrumentation and PowerShell.'
  - 'How attack software ''S0363: Empire'' uses attack technique ''T1558.001: Golden Ticket'':

    Empire can leverage its implementation of Mimikatz to obtain and use golden tickets.'
  - 'Description of attack group ''G0024: Putter Panda'':

    Putter Panda is a Chinese threat group that has been attributed to Unit 61486 of the 12th Bureau of the PLA’s 3rd General
    Staff Department (GSD).'
- source_sentence: 'Represent this sentence for searching relevant passages: What is Carberp?'
  sentences:
  - 'How attack software ''S0556: Pay2Key'' uses attack technique ''T1573.002: Asymmetric Cryptography'':

    Pay2Key has used RSA encrypted communications with C2.'
  - 'How attack software ''S0165: OSInfo'' uses attack technique ''T1016: System Network Configuration Discovery'':

    OSInfo discovers the current domain information.'
  - 'Description of attack software ''S0484: Carberp'':

    Carberp is a credential and information stealing malware that has been active since at least 2009. Carberp''s source code
    was leaked online in 2013, and subsequently used as the foundation for the Carbanak backdoor.'
- source_sentence: 'Represent this sentence for searching relevant passages: What attack techniques are used by software ''S0144:
    ChChes''?'
  sentences:
  - 'How attack software ''S1119: LIGHTWIRE'' uses attack technique ''T1573.001: Symmetric Cryptography'':

    LIGHTWIRE can RC4 encrypt C2 commands.'
  - 'The attack techniques used by software ''S0144: ChChes'' are: ''T1036.005: Match Legitimate Name or Location'', ''T1057:
    Process Discovery'', ''T1071.001: Web Protocols'', ''T1082: System Information Discovery'', ''T1083: File and Directory
    Discovery'', ''T1105: Ingress Tool Transfer'', ''T1132.001: Standard Encoding'', ''T1547.001: Registry Run Keys / Startup
    Folder'', ''T1553.002: Code Signing'', ''T1555.003: Credentials from Web Browsers'', ''T1562.001: Disable or Modify Tools'',
    ''T1573.001: Symmetric Cryptography'''
  - 'How attack software ''S0649: SMOKEDHAM'' uses attack technique ''T1204.001: Malicious Link'':

    SMOKEDHAM has relied upon users clicking on a malicious link delivered through phishing.'
model-index: null
config:
  architectures:
  - BertModel
  model_type: bert
  tokenizer_config:
    cls_token: '[CLS]'
    mask_token: '[MASK]'
    pad_token: '[PAD]'
    sep_token: '[SEP]'
    unk_token: '[UNK]'
cardData:
  tags:
  - sentence-transformers
  - sentence-similarity
  - feature-extraction
  - generated_from_trainer
  - dataset_size:22801
  - loss:MultipleNegativesRankingLoss
  base_model: BAAI/bge-base-en-v1.5
  widget:
  - source_sentence: 'Represent this sentence for searching relevant passages: How does attack software ''S0030: Carbanak''
      use attack technique ''T1027: Obfuscated Files or Information''?'
    sentences:
    - 'How attack software ''S0363: Empire'' uses attack technique ''T1558.002: Silver Ticket'':

      Empire can leverage its implementation of Mimikatz to obtain and use silver tickets.'
    - 'How attack software ''S0664: Pandora'' uses attack technique ''T1027: Obfuscated Files or Information'':

      Pandora has the ability to compress stings with QuickLZ.'
    - 'How attack software ''S0030: Carbanak'' uses attack technique ''T1027: Obfuscated Files or Information'':

      Carbanak encrypts strings to make analysis more difficult.'
  - source_sentence: 'Represent this sentence for searching relevant passages: What type of organizations has MegaCortex mainly
      targeted?'
    sentences:
    - 'How attack software ''S0534: Bazar'' uses attack technique ''T1070.009: Clear Persistence'':

      Bazar''s loader can delete scheduled tasks created by a previous instance of the malware.'
    - 'Description of attack software ''S0576: MegaCortex'':

      MegaCortex is ransomware that first appeared in May 2019.  MegaCortex has mainly targeted industrial organizations.'
    - 'How attack software ''S0093: Backdoor.Oldrea'' uses attack technique ''T1057: Process Discovery'':

      Backdoor.Oldrea collects information about running processes.'
  - source_sentence: 'Represent this sentence for searching relevant passages: Describe attack group ''G0024: Putter Panda'''
    sentences:
    - 'How data component ''Script Execution'' can be used to detect attack technique ''T1005: Data from Local System'':

      Monitor for any attempts to enable scripts running on a system would be considered suspicious. If scripts are not commonly
      used on a system, but enabled, scripts running out of cycle from patching or other administrator functions are suspicious.
      Scripts should be captured from the file system when possible to determine their actions and intent.  Data may also
      be acquired through Windows system management tools such as Windows Management Instrumentation and PowerShell.'
    - 'How attack software ''S0363: Empire'' uses attack technique ''T1558.001: Golden Ticket'':

      Empire can leverage its implementation of Mimikatz to obtain and use golden tickets.'
    - 'Description of attack group ''G0024: Putter Panda'':

      Putter Panda is a Chinese threat group that has been attributed to Unit 61486 of the 12th Bureau of the PLA’s 3rd General
      Staff Department (GSD).'
  - source_sentence: 'Represent this sentence for searching relevant passages: What is Carberp?'
    sentences:
    - 'How attack software ''S0556: Pay2Key'' uses attack technique ''T1573.002: Asymmetric Cryptography'':

      Pay2Key has used RSA encrypted communications with C2.'
    - 'How attack software ''S0165: OSInfo'' uses attack technique ''T1016: System Network Configuration Discovery'':

      OSInfo discovers the current domain information.'
    - 'Description of attack software ''S0484: Carberp'':

      Carberp is a credential and information stealing malware that has been active since at least 2009. Carberp''s source
      code was leaked online in 2013, and subsequently used as the foundation for the Carbanak backdoor.'
  - source_sentence: 'Represent this sentence for searching relevant passages: What attack techniques are used by software
      ''S0144: ChChes''?'
    sentences:
    - 'How attack software ''S1119: LIGHTWIRE'' uses attack technique ''T1573.001: Symmetric Cryptography'':

      LIGHTWIRE can RC4 encrypt C2 commands.'
    - 'The attack techniques used by software ''S0144: ChChes'' are: ''T1036.005: Match Legitimate Name or Location'', ''T1057:
      Process Discovery'', ''T1071.001: Web Protocols'', ''T1082: System Information Discovery'', ''T1083: File and Directory
      Discovery'', ''T1105: Ingress Tool Transfer'', ''T1132.001: Standard Encoding'', ''T1547.001: Registry Run Keys / Startup
      Folder'', ''T1553.002: Code Signing'', ''T1555.003: Credentials from Web Browsers'', ''T1562.001: Disable or Modify
      Tools'', ''T1573.001: Symmetric Cryptography'''
    - 'How attack software ''S0649: SMOKEDHAM'' uses attack technique ''T1204.001: Malicious Link'':

      SMOKEDHAM has relied upon users clicking on a malicious link delivered through phishing.'
  pipeline_tag: sentence-similarity
  library_name: sentence-transformers
siblings:
- rfilename: .gitattributes
- rfilename: 1_Pooling/config.json
- rfilename: README.md
- rfilename: config.json
- rfilename: config_sentence_transformers.json
- rfilename: model.safetensors
- rfilename: modules.json
- rfilename: sentence_bert_config.json
- rfilename: tokenizer.json
- rfilename: tokenizer_config.json
spaces: []
createdAt: '2026-04-27T07:49:32.000Z'
safetensors:
  parameters:
    F32: 109482240
  total: 109482240
usedStorage: 437951304
downloadsAllTime: 61
inferenceProviderMapping: {}
configJson:
  add_cross_attention: false
  architectures:
  - BertModel
  attention_probs_dropout_prob: 0.1
  bos_token_id: null
  classifier_dropout: null
  dtype: float32
  eos_token_id: null
  gradient_checkpointing: false
  hidden_act: gelu
  hidden_dropout_prob: 0.1
  hidden_size: 768
  id2label:
    '0': LABEL_0
  initializer_range: 0.02
  intermediate_size: 3072
  is_decoder: false
  label2id:
    LABEL_0: 0
  layer_norm_eps: 1.0e-12
  max_position_embeddings: 512
  model_type: bert
  num_attention_heads: 12
  num_hidden_layers: 12
  pad_token_id: 0
  position_embedding_type: absolute
  tie_word_embeddings: true
  transformers_version: 5.0.0
  type_vocab_size: 2
  use_cache: true
  vocab_size: 30522
onnxFiles: {}
safetensorsSizeMb: 438.0
approxSizeMb: 437.9
fetchedAt: '2026-06-09T06:55:53.000Z'
